# s6/services/sshd/run
#!/usr/bin/execlineb -P
/usr/sbin/sshd -D -E /var/log/sshd.log
