# s6/services/nginx/run
#!/usr/bin/execlineb -P
/usr/sbin/nginx -c /etc/services.d/nginx/nginx.conf
